# crontab
0,30 * * * * /usr/bin/env bash $HOME/BoostStats/automate BoostStats scraper_ci frontend_ci >> $HOME/BoostStats/automate.log 2>&1
